# prometheus/log-monitor.timer
[Unit]
Description= Run logging tasks hourly

[Timer]
OnBootSec=1min
OnUnitActiveSec=1m
Unit=log-monitor.target

[Install]
WantedBy=timers.target
